# docs/deployment/russell-sentinel.service
# SPDX-License-Identifier: MIT OR Apache-2.0
# NOTE: The canonical systemd units live in packaging/systemd/.
# This file is retained for the docs/deployment/install.sh reference.
[Unit]
Description=Russell Sentinel — Host probe collection
Documentation=https://github.com/Replicant-Partners/russell
After=network.target

[Service]
Type=oneshot
TimeoutStartSec=120
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/state/harness %h/.local/share/harness
PrivateTmp=true
ExecStart=%h/.local/bin/russell sentinel-once
StandardOutput=journal
StandardError=journal
SyslogIdentifier=russell-sentinel

[Install]
WantedBy=default.target
