# packaging/systemd/kask-gateway.service
[Unit]
Description=Kask MCP Gateway (stack-api)
Documentation=https://kask.ai
After=network.target
PartOf=default.target

[Service]
Type=exec
ExecStart=%h/Clones/kask/target/debug/stack-api
Restart=on-failure
RestartSec=5s
TimeoutStopSec=10s

# Environment
Environment=RUST_LOG=info
Environment=KASK_BIND_ADDR=127.0.0.1:18100

# Security hardening
NoNewPrivs=true
ProtectHome=true
ProtectSystem=true
PrivateTmp=true

[Install]
WantedBy=default.target