# packaging/systemd/kask-token-rotate.service
[Unit]
Description=Rotate Russell Kask MCP Token
Documentation=file:///home/mdz-axolotl/Clones/russell/docs/operations/KASK_TOKEN_ROTATION.md

[Service]
Type=oneshot
ExecStart=%h/.local/bin/rotate-russell-token.sh
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivs=true
ProtectHome=true
ProtectSystem=true
PrivateTmp=true
