# packaging/systemd/russell-okapi.service
# SPDX-License-Identifier: MIT OR Apache-2.0
# Run once per Okapi probe cycle. Installed to ~/.config/systemd/user/.

[Unit]
Description=Russell Okapi Probe — observe the inference engine
OnFailure=russell-failure@%n.service

[Service]
Type=oneshot
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/state/harness %h/.local/share/harness
PrivateTmp=yes
NoNewPrivileges=yes
EnvironmentFile=%h/.config/harness/russell.env
ExecStart=%h/.local/bin/russell okapi-probe --auto-apply ${RUSSELL_OKAPI_DEFAULT_MODEL:+--default-model $RUSSELL_OKAPI_DEFAULT_MODEL}
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target