# packaging/systemd/russell-sentinel.timer
# SPDX-License-Identifier: MIT OR Apache-2.0
# Fires the Sentinel every 5 minutes, catching up after sleep.

[Unit]
Description=Russell Sentinel — 5-minute cadence
Documentation=file://%h/Clones/russell/docs/specifications/MVP_SPEC.md

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
RandomizedDelaySec=15s
Persistent=true
Unit=russell-sentinel.service

[Install]
WantedBy=timers.target
