# packaging/systemd/russell-digest.service
# SPDX-License-Identifier: MIT OR Apache-2.0
# Weekly Markdown digest of Russell's observations.

[Unit]
Description=Russell weekly digest
OnFailure=russell-failure@%p.service

[Service]
Type=oneshot
Nice=19
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/state/harness
PrivateTmp=yes
NoNewPrivileges=yes
ExecStartPre=/bin/mkdir -p %h/.local/state/harness/digest
ExecStart=/bin/sh -c '%h/.local/bin/russell digest --since-hours 168 > %h/.local/state/harness/digest/$(date -u +%%Y-W%%V).md'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
