# packaging/systemd/russell-okapi.timer
# SPDX-License-Identifier: MIT OR Apache-2.0
# Fires the Okapi probe every 5 minutes, offset from Sentinel.

[Unit]
Description=Russell Okapi Probe — 5-minute cadence

[Timer]
OnBootSec=4min30s
OnUnitActiveSec=5min
AccuracySec=30s
RandomizedDelaySec=15s
Persistent=true
Unit=russell-okapi.service

[Install]
WantedBy=timers.target