# packaging/systemd/russell-digest.timer
# SPDX-License-Identifier: MIT OR Apache-2.0
# Sunday 09:00 local — the operator's weekly review.

[Unit]
Description=Russell weekly digest timer
Documentation=file://%h/Clones/russell/docs/specifications/MVP_SPEC.md

[Timer]
OnCalendar=Sun 09:00
Persistent=true
RandomizedDelaySec=10m
Unit=russell-digest.service

[Install]
WantedBy=timers.target
